# server/layout-api.service
[Unit]
Description=PP-DocLayout API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/xoaghim/server
Environment="PATH=/root/layout_api_venv/bin"
ExecStart=/root/layout_api_venv/bin/uvicorn layout_api_server:app --host 0.0.0.0 --port 8765
Restart=always
RestartSec=10

# GPU support
Environment="CUDA_VISIBLE_DEVICES=0"

[Install]
WantedBy=multi-user.target
